# service/claude-agent.service
[Unit]
Description=Claude Agent System
After=network.target

[Service]
Type=simple
User=claude
Group=claude
WorkingDirectory=/home/claude/claude-agent-system
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/home/claude/.local/bin"
ExecStart=/bin/bash -c 'if [ -f /home/claude/claude-agent-system/orchestrator.pid ]; then exec tail -f /dev/null; fi'
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
